# appconfig/i3/systemd/battery_notification.service
[Unit]
Description=Sends notification about battery status

[Service]
# %h means ~
ExecStart=%h/.scripts/battery_notification.sh
Type=oneshot
